# src/eye-tracker/digitizer/hid-digitizer.service
[Unit]
Description=HID Digitizer Controller Service
Documentation=https://github.com/dimitarkyurtov/touchless-hid-digitizer
After=usb-gadget.service network.target
Requires=usb-gadget.service
Wants=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/touchless-hid-digitizer/eye-tracker/digitizer
ExecStart=/usr/bin/python3 /opt/touchless-hid-digitizer/eye-tracker/digitizer/main.py
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/log

[Install]
WantedBy=multi-user.target
